# scripts/chat-service.service
[Unit]
Description=Chat Demo Service (SparkChain)
After=network.target asr-service.service
Wants=network.target
Requires=asr-service.service

[Service]
Type=simple
User=bxi
Group=bxi
WorkingDirectory=/home/bxi/M2_SDK/星火SDK/SparkChain_Linux_SDK_2.0.0_rc1/llm/online_llm/chat

# 使用包装脚本启动，解决库兼容性问题
ExecStart=/home/bxi/M2_SDK/scripts/chat_wrapper.sh

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 增加启动超时时间
TimeoutStartSec=30
TimeoutStopSec=10

# 确保日志目录和脚本权限，等待ASR服务启动
ExecStartPre=/bin/mkdir -p /home/bxi/M2_SDK/log
ExecStartPre=/bin/chmod +x /home/bxi/M2_SDK/scripts/chat_wrapper.sh
ExecStartPre=/bin/sleep 10

[Install]
WantedBy=multi-user.target
